# stow.d/profile/.config/profile.d/fzf.sh
# Base colors from Kanagawa palette
export FZF_DEFAULT_OPTS="
  --color=fg:#DCD7BA,bg:#1F1F28,hl:#957FB8
  --color=fg+:#C8C093,bg+:#2A2A37,hl+:#D27E99
  --color=info:#7E9CD8,prompt:#76946A,pointer:#DCA561
  --color=marker:#98BB6C,spinner:#7FB4CA,header:#FFA066
  --color=border:#54546D,gutter:#1F1F28
  --border=rounded
  --prompt='❯ '
  --pointer='➜'
  --marker='✓'
"
